# kali-fs/var/lib/dpkg/info/polkitd.list
/.
/etc
/etc/polkit-1
/etc/polkit-1/rules.d
/etc/xml
/usr
/usr/bin
/usr/bin/pkaction
/usr/bin/pkcheck
/usr/bin/pkttyagent
/usr/lib
/usr/lib/pam.d
/usr/lib/pam.d/polkit-1
/usr/lib/policykit-1
/usr/lib/polkit-1
/usr/lib/polkit-1/polkit-agent-helper-1
/usr/lib/polkit-1/polkitd
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/polkit.service
/usr/lib/sysusers.d
/usr/lib/sysusers.d/polkitd.conf
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/polkitd.conf
/usr/libexec
/usr/share
/usr/share/bug
/usr/share/bug/polkitd
/usr/share/bug/polkitd/control
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.PolicyKit1.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.freedesktop.PolicyKit1.conf
/usr/share/doc
/usr/share/doc/polkitd
/usr/share/doc/polkitd/NEWS.Debian.gz
/usr/share/doc/polkitd/NEWS.md.gz
/usr/share/doc/polkitd/README.md
/usr/share/doc/polkitd/changelog.Debian.gz
/usr/share/doc/polkitd/copyright
/usr/share/doc/polkitd/examples
/usr/share/doc/polkitd/examples/50-local-allow.rules
/usr/share/doc/polkitd/examples/50-local-deny.rules
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/polkitd
/usr/share/locale
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/polkit-1.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/polkit-1.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/polkit-1.mo
/usr/share/locale/hr
/usr/share/locale/hr/LC_MESSAGES
/usr/share/locale/hr/LC_MESSAGES/polkit-1.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/polkit-1.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/polkit-1.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/polkit-1.mo
/usr/share/locale/ka
/usr/share/locale/ka/LC_MESSAGES
/usr/share/locale/ka/LC_MESSAGES/polkit-1.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/polkit-1.mo
/usr/share/locale/nn
/usr/share/locale/nn/LC_MESSAGES
/usr/share/locale/nn/LC_MESSAGES/polkit-1.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/polkit-1.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/polkit-1.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/polkit-1.mo
/usr/share/locale/ro
/usr/share/locale/ro/LC_MESSAGES
/usr/share/locale/ro/LC_MESSAGES/polkit-1.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/polkit-1.mo
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/polkit-1.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/polkit-1.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/tr/LC_MESSAGES/polkit-1.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/polkit-1.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/polkit-1.mo
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/polkit-1.mo
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pkaction.1.gz
/usr/share/man/man1/pkcheck.1.gz
/usr/share/man/man1/pkttyagent.1.gz
/usr/share/man/man8
/usr/share/man/man8/polkit.8.gz
/usr/share/man/man8/polkitd.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.policykit.policy
/usr/share/polkit-1/policyconfig-1.dtd
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/50-default.rules
/usr/share/xml
/usr/share/xml/polkit-1
/usr/share/xml/polkit-1/catalog.xml
/var
/var/lib
/var/lib/polkit-1
/usr/lib/policykit-1/polkit-agent-helper-1
/usr/libexec/polkit-agent-helper-1
